# bot.conf
[program:cafebot2]
directory=/home/ubuntu/KafeBot2
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu", USER="ubuntu"
stderr_logfile=/home/ubuntu/KafeBot2/logfile_err.log
stdout_logfile=/home/ubuntu/KafeBot2/logfile.log
